# services/vpn.service
# Contents of /etc/systemd/system/vpn.service
[Unit]
Description=Starts VPN service
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=SERVICES_FOLDER/scripts/vpn-connect.sh start
ExecStop=SERVICES_FOLDER/scripts/vpn-connect.sh stop

[Install]
WantedBy=multi-user.target